# config/supervisor/purbeurre-gunicorn.conf
[program:purbeurre]
directory = /home/ubuntu/oc_dapython_pr10/
command = /home/ubuntu/.local/bin/pipenv run gunicorn --chdir purbeurre purbeurre.wsgi:application
autostart = true
autorestart = true
user = ubuntu
environment = DJANGO_SETTINGS_MODULE='purbeurre.settings.aws_settings'

[supervisord]
environment = LC_ALL='en_US.UTF-8',LANG='en_US.UTF-8'