# config/supervisor/purbeurre-gunicorn.conf
[program:purbeurre]
directory=/home/ubuntu/oc_dapython_pr10/purbeurre/
command=pipenv run gunicorn --chdir purbeurre purbeurre.wsgi
autostart = true
autorestart = true